# systemd/dns-healthcheck.timer
[Unit]
Description=Zelira — DNS Health Check Timer (every 2 min)

[Timer]
OnBootSec=60
OnUnitActiveSec=120
AccuracySec=10

[Install]
WantedBy=timers.target
